# base/templates/default/clean_tmp.cron
0 */4 * * * root find /tmp -maxdepth 1 -type f -amin +720 -exec rm -f {} \;